# deploy/whoisfreaks-mcp.service
[Unit]
Description=WhoisFreaks MCP Server
Documentation=https://github.com/whoisfreaks/whoisfreaks-mcp-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=whoisfreaks
Group=whoisfreaks

# JAR location — adjust if you installed to a different path
ExecStart=/usr/bin/java -jar /opt/whoisfreaks-mcp/whoisfreaks-mcp-server-1.0.0.jar

# Load secrets from a separate file (chmod 600, owned by root)
# File contains a single line: WHOISFREAKS_API_KEY=your-key-here
EnvironmentFile=/etc/whoisfreaks-mcp/env

# Restart automatically on crash (but not on clean exit)
Restart=on-failure
RestartSec=5

# Give the JVM enough file descriptors
LimitNOFILE=65536

# stdout is reserved for the MCP stdio protocol — send it to /dev/null.
# All application logs (stderr) go to journald: journalctl -u whoisfreaks-mcp -f
StandardOutput=null
StandardError=journal
SyslogIdentifier=whoisfreaks-mcp

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/whoisfreaks-mcp

[Install]
WantedBy=multi-user.target
